# kodi/.kodi/userdata/epg/epg.cron
0 * * * * /storage/.kodi/userdata/epg/epg
